# system-files/usr/lib/systemd/system/repair-accounts.service
[Unit]
Description=Repair orphaned passwd/group shadow database entries
DefaultDependencies=no
After=systemd-remount-fs.service
Before=systemd-sysusers.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/repair-accounts.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
